# 시스템프로그래밍실습 | System Programming Lab (SWE2024)/lab workout/[SPL] week6/simple-cron.service
[Unit]
Description=Simple Cron Daemon
After=network.target

[Service]
Type=forking
ExecStart=/SPL/week6/simple_cron
RemainAfterExit=yes

[Install]
WantedBy=default.target